# DEBIAN/postrm
#!/bin/bash
systemctl stop contextNetServer.service
systemctl disable contextNetServer.service
systemctl daemon-reload
