# DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable contextNetServer.service
systemctl start contextNetServer.service