# DEBIAN/postinst
#!/bin/bash
systemctl enable contextNetServer.service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl start contextNetServer.service
fi
