# etc/profile.d/opensplice.sh
../../opt/contextNetServer/libs/OpenSplice/opensplice.sh